# infra/systemd/mstechnics-vnnox-unresolved.service
[Unit]
Description=MsTechnics VNNOX unresolved alarm notification check
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/mstechnics
EnvironmentFile=/opt/mstechnics/Config/.env
ExecStart=/opt/mstechnics/.venv/bin/python manage.py check_unresolved_alarms

