# infra/systemd/mstechnics-vnnox-unresolved.service
[Unit]
Description=MsTechnics VNNOX unresolved alarm notification check
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/DisplayControl/MsTechnics
EnvironmentFile=/root/DisplayControl/MsTechnics/Config/.env
ExecStart=/root/DisplayControl/venv-c84d3816/bin/python manage.py check_unresolved_alarms
