# infra/systemd/mstechnics-db-backup.service
[Unit]
Description=MsTechnics PostgreSQL backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/mstechnics
EnvironmentFile=/opt/mstechnics/Config/.env
ExecStart=/opt/mstechnics/scripts/backup-db.sh
