# infra/systemd/mstechnics-gunicorn.service
[Unit]
Description=MsTechnics gunicorn with gevent workers for SSE
After=network.target postgresql.service redis-server.service

[Service]
WorkingDirectory=/root/DisplayControl/MsTechnics
EnvironmentFile=/root/DisplayControl/MsTechnics/Config/.env
ExecStart=/root/DisplayControl/venv-c84d3816/bin/gunicorn \
    --workers 3 \
    --worker-class gevent \
    --worker-connections 1000 \
    --bind unix:/run/gunicorn.sock \
    --timeout 120 \
    project_config.wsgi:application
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
