# infra/systemd/mstechnics-db-backup.timer
[Unit]
Description=Run MsTechnics PostgreSQL backup once per day

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=mstechnics-db-backup.service

[Install]
WantedBy=timers.target
