# infra/systemd/mstechnics-db-backup.service
[Unit]
Description=MsTechnics PostgreSQL backup
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/DisplayControl/MsTechnics
EnvironmentFile=/root/DisplayControl/MsTechnics/Config/.env
Environment=BACKUP_DIR=/var/backups/mstechnics/scheduled
ExecStart=/bin/bash /root/DisplayControl/MsTechnics/scripts/backup-db.sh
